# upstart_template.conf
description "%APP_NAME%"
author      "Sprocket http://www.rafflecopter.com"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/app-data"
    export NODE_ENV=production

    chdir /var/builds/live/%APP_NAME%/%APP_NAME% 

    exec node /var/builds/live/%APP_NAME%/%APP_NAME%/service.js >> /var/log/%APP_NAME%-live.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
